# cron/huobi_price.cron
*/1 * * * * bash /mnt/share/gitStuff/recipes/cron/huobi_price_recorder.sh >>/dev/null 2>&1
*/2 * * * * bash /mnt/share/gitStuff/recipes/cron/huobi_price_alert.sh >>/dev/null 2>&1
*/5 * * * * bash /mnt/share/gitStuff/recipes/cron/chinese_btc_difference.sh >>/dev/null 2>&1
